# launchd/com.xsiam.threat-intel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xsiam.threat-intel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/dhruvsethi/Downloads/Personal Projects/claude-xsiam-data-aggregator/pipeline.py</string>
    </array>

    <!-- Fires daily at 2:00 PM Singapore time (06:00 UTC) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/dhruvsethi/Downloads/Personal Projects/claude-xsiam-data-aggregator</string>

    <key>StandardOutPath</key>
    <string>/tmp/xsiam-threat-intel.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/xsiam-threat-intel-error.log</string>

    <!-- Run once immediately if the scheduled time was missed (e.g. laptop was asleep) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
